# crontab
30 * * * * rclone sync gdrive:/ /mnt/data/rclone/gdrive